# debian/prerm
#!/bin/sh -e
#DEBHELPER#

invoke-rc.d loggerhead stop
