# debian/prerm
#!/bin/sh -e

invoke-rc.d loggerhead stop || exit 0
